# systemd/spns.target
[Unit]
Description=Session Push Notification Server services
